# __back/AlgoritmPlusBot.service
[Unit]
Description=AlgoritmPlus Telegram Bot
After=network-online.target systemd-resolved.service
Wants=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/projects/AlgoritmPlusBot/bot

ExecStart=/var/projects/AlgoritmPlusBot/bot/.venv/bin/python3 /var/projects/AlgoritmPlusBot/bot/run.py

Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="TZ=Europe/Moscow"

ReadWritePaths=/var/projects/AlgoritmPlusBot/bot

StandardOutput=journal
StandardError=journal
SyslogIdentifier=AlgoritmPlusBot

[Install]
WantedBy=multi-user.target
